# script/starman_centos
#!/bin/bash
 
<<<<<<< HEAD
# subman - this script starts and stops the subman daemon
#
# chkconfig:   - 85 15
# description: Subman 
# processname: subman
# pidfile :    /var/run/starman.pid
# www file:    /opt/subman
=======
# chargemonk - this script starts and stops the chargemonk daemon
#
# chkconfig:   - 85 15
# description: Chargemonk 
# processname: chargemonk
# pidfile :    /var/run/starman.pid
# www file:    /opt/chargemonk
>>>>>>> master
 
 
# Source function library.
. /etc/rc.d/init.d/functions
 
 
export plackup="/usr/local/bin/plackup"
export starman="/usr/local/bin/starman"
<<<<<<< HEAD
export myapp="subman"
export myapp_path="/opt/subman"
export pidfile="/var/run/subman.pid"
=======
export myapp="chargemonk"
export myapp_path="/opt/chargemonk"
export pidfile="/var/run/chargemonk.pid"
>>>>>>> master
export app_port="3000"
 
start() {
	if [ -f "${pidfile}" ] 
	then  
		echo -e "Process already started!"
	else
	        echo -n $"Starting ${myapp}: "
		$plackup -I${myapp_path}/lib  --user pearlbee --listen :${app_port} --listen /tmp/${myapp}.sock -D --pid ${pidfile} -E deployment -s Starman --error-log /var/log/nginx/${myapp}.log -a ${myapp_path}/${myapp}.psgi 
	        RETVAL=$?
	        echo  
	        [ $RETVAL = 1 ]
	        return $RETVAL
	fi
}
 
restart() {
        echo -n $"Restarting ${myapp}: "
	stop
	start
}
 
stop() {
        echo -n $"Stopping ${myapp}: "
        killproc -p ${pidfile} -d 10 $starman
        RETVAL=$?
        echo
        [ $RETVAL = 0 ] && rm -f ${pidfile}
}
# See how we were called.
case "$1" in
  start)
        start
        ;;
  restart)
        restart
         ;;
  stop)
        stop
        ;;
  *)
        echo $"Usage: ${myapp} {start|restart|stop}"
        exit 1
esac
 
exit $RETVAL
